# src/setup/misc/backup.cron
# run backup script every day at 1:00
00 01 * * * /opt/hilel14/archie/beeri/bin/backup.sh
